# launchd/agency.athos.hub.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>agency.athos.hub</string>
  <key>WorkingDirectory</key>
  <string>/Users/clem/Sites/athos</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/clem/Sites/athos/venv/bin/python</string>
    <string>/Users/clem/Sites/athos/voice/server.py</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>SessionCreate</key>
  <true/>
  <key>LimitLoadToSessionType</key>
  <string>Aqua</string>
  <key>StandardOutPath</key>
  <string>/Users/clem/Sites/athos/logs/athos_launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/clem/Sites/athos/logs/athos_launchd.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/clem</string>
    <key>USER</key>
    <string>clem</string>
    <key>DRIVE_PATH</key>
    <string>/Users/clem/Sites/athos/memory</string>
    <key>ATHOS_PATH</key>
    <string>/Users/clem/Sites/athos</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/Users/clem/.vscode/extensions/openai.chatgpt-26.5513.21555-darwin-x64/bin/macos-x86_64:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
